# docs/systemd/relay2slack.service.d/env.conf
[Service]
Environment=HTTP_PROXY=http://your.proxy.tld:3128
Environment=HTTPS_PROXY=http://your.proxy.tld:3128
Environment=SLACK_TOKEN=/Y0UR/70K3N/G0E$H3R3
