# docs/systemd/relay2slack.service.d/proxy.conf
[Service]
Environment=HTTP_PROXY=http://your.proxy.tld:3128
Environment=HTTPS_PROXY=http://your.proxy.tld:3128
